# znf.service
[Unit]
Description=ZNF Web Application
After=network.target redis.service
Wants=redis.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/usr/share/znf
Environment=ZNF_CONFIG=/usr/share/znf/config.ini
ExecStart=/usr/share/znf/.venv/bin/gunicorn --config gunicorn.conf.py server:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=znf

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/usr/share/znf/logs /usr/share/znf/server.sock /mnt/files/znf

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
